# host/install/midex-upload@.service
# Oneshot that uploads the class-compliant MIDI firmware into a loader-PID MIDEX.
#
# Triggered (not enabled) by udev: 99-midex-class.rules sets SYSTEMD_WANTS on the
# device, so systemd starts one instance per appearing device. The instance name
# (%i) is the USB device's kernel name (e.g. 1-2); the dispatch wrapper uses it to
# read the device's product ID from sysfs and pick the matching firmware image
# (r1 for 0x1000, r2 for 0x1010), so auto-upload is revision-aware.
#
# To switch back to the stock vendor firmware, mask this unit and let the
# in-tree snd_usb_midex driver claim the device instead:
#   host/midex-fwctl.sh stock     (masks this unit, unblacklists snd_usb_midex)

[Unit]
Description=Upload class-compliant MIDI firmware to MIDEX (%i)

[Service]
Type=oneshot
ExecStart=/usr/local/bin/midex-upload-dispatch.sh %i
# -ENODEV on the final CPUCS=0 is expected (the device re-enumerates mid-request)
# and the uploader already treats it as success, so a clean exit is correct here.
